# note-daemon.service
[Unit]
Description=NetNotes Secure Daemon
After=network.target

[Service]
Type=simple
User=root
Group=netnotes
RuntimeDirectory=netnotes
RuntimeDirectoryMode=0755

ExecStart=/usr/local/bin/note-daemon
Restart=no
ExecStop=/bin/kill -TERM 

# Optional: capture logs in journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target